# docker/supervisord.conf
; ============================================================
; MES Pangan Masa Depan — Supervisord Configuration
; ============================================================
; Manages 4 processes: PostgreSQL, Node.js Backend, Nginx, ML Service
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ---- PostgreSQL 16 ----
; Priority 10 = starts first, stops last
[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; ---- Node.js Backend API ----
; Priority 20 = starts after PostgreSQL
[program:backend]
command=node /app/dist/index.js
directory=/app
user=root
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=5
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=
    NODE_ENV="%(ENV_NODE_ENV)s",
    PORT="%(ENV_PORT)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    ML_SERVICE_URL="%(ENV_ML_SERVICE_URL)s",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    SEED_SECRET="%(ENV_SEED_SECRET)s",
    FRONTEND_URL="%(ENV_FRONTEND_URL)s"

; ---- Nginx (Frontend + Reverse Proxy) ----
; Priority 30 = starts after Backend
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
startsecs=2
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; ---- Python ML Service ----
; Priority 20 = starts alongside Backend
[program:ml-service]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app/ml-service
user=root
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/ml-service.log
stderr_logfile=/var/log/supervisor/ml-service-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
